# dist/cloudflare-railgun/files/supervisord/apache.conf
[program:apache]
command = /bin/start-apache
autostart = true
autorestart = true
stdout_events_enabled = true 
stderr_events_enabled = true
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 0
stderr_logfile_maxbytes = 1MB
stderr_logfile_backups = 0

[program:apache-access-log]
command = tail -f -n0 /var/log/apache2/access.log
autostart = true
autorestart = true
stdout_events_enabled = true 
stderr_events_enabled = true
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 0
stderr_logfile_maxbytes = 1MB
stderr_logfile_backups = 0

[program:apache-error-log]
command = tail -f -n0 /var/log/apache2/error.log
autostart = true
autorestart = true
stdout_events_enabled = true 
stderr_events_enabled = true
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 0
stderr_logfile_maxbytes = 1MB
stderr_logfile_backups = 0