# dist/mongodb/files/supervisord/php-fpm70.conf
[program:php-fpm70]
command = /usr/sbin/php-fpm7.0 -c /etc/php/7.0/fpm/php-fpm.conf -F
autostart = true
autorestart = true
stdout_events_enabled = true 
stderr_events_enabled = true
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 0
stderr_logfile_maxbytes = 1MB
stderr_logfile_backups = 0
