# dist/apache-php56/files/supervisord/gitlab.conf
[program:unicorn_rails]
user=git
command=/usr/local/bin/bundle exec unicorn_rails -c /home/git/gitlab/config/unicorn.rb -E production
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/home/git/gitlab/log/unicorn.stdout.log
stderr_logfile=/home/git/gitlab/log/unicorn.stderr.log
directory=/home/git/gitlab
environment=RAILS_ENV="production",HOME="/home/git"
stdout_events_enabled = true 
stderr_events_enabled = true
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 0
stderr_logfile_maxbytes = 1MB
stderr_logfile_backups = 0

[program:gitlab_logs]
command = tail -f -n0 /home/git/gitlab/logs/production.log
autostart = true
autorestart = true
stdout_events_enabled = true 
stderr_events_enabled = true
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 0
stderr_logfile_maxbytes = 1MB
stderr_logfile_backups = 0

[program:sidekiq]
user=git
command=/usr/local/bin/bundle exec sidekiq -q post_receive -q mailers -q archive_repo -q system_hook -q project_web_hook -q gitlab_shell -q incoming_email -q runner -q common -q default -e production -P /home/git/gitlab/tmp/pids/sidekiq.pid
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/home/git/gitlab/log/sidekiq-stdout.log
stderr_logfile=/home/git/gitlab/log/sidekiq-stderr.log
directory=/home/git/gitlab
environment=RAILS_ENV="production",HOME="/home/git"
stdout_events_enabled = true 
stderr_events_enabled = true
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 0
stderr_logfile_maxbytes = 1MB
stderr_logfile_backups = 0

[program:gitlab_workhorse]
user=git
command=/home/git/gitlab-workhorse/gitlab-workhorse -listenUmask 0 -listenNetwork unix -listenAddr /home/git/gitlab/tmp/sockets/gitlab-workhorse.socket -authBackend http://127.0.0.1:8080 -authSocket /home/git/gitlab/tmp/sockets/gitlab.socket -documentRoot /home/git/gitlab/public
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/home/git/gitlab/log/gitlab-workhorse-stdout.log
stderr_logfile=/home/git/gitlab/log/gitlab-workhorse-stderr.log
directory=/home/git/gitlab
environment=RAILS_ENV="production",HOME="/home/git"
stdout_events_enabled = true 
stderr_events_enabled = true
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 0
stderr_logfile_maxbytes = 1MB
stderr_logfile_backups = 0

[program:mail_room]
user=git
command=/usr/local/bin/bundle exec mail_room -q -c /home/git/gitlab/config/mail_room.yml
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/home/git/gitlab/log/mail_room-stdout.log
stderr_logfile=/home/git/gitlab/log/mail_room-stderr.log
directory=/home/git/gitlab
environment=RAILS_ENV="production",HOME="/home/git"
stdout_events_enabled = true 
stderr_events_enabled = true
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 0
stderr_logfile_maxbytes = 1MB
stderr_logfile_backups = 0

