# images/nginx-php/shared/php-fpm/rootfs/etc/service/php-fpm/run
#!/bin/sh
exec /usr/sbin/php-fpm$PHP_VERSION -c /etc/php$PHP_VERSION/php-fpm.conf -F
