# files/wordpress-nginx-w3tc/crontab
*/5 * * * * /usr/bin/sudo -u www-data /usr/bin/php /var/www/wp-cron.php
